# deploy-bundle/systemd/blakestream-mpos-backup.service
[Unit]
Description=Blakestream-MPOS backup snapshot
After=mariadb.service docker.service
Wants=mariadb.service

[Service]
Type=oneshot
User=root
# backup.sh sources /opt/blakestream-mpos/.deploy.env itself, so no
# need to pre-source from the unit. TimeoutStartSec gives the wallet
# RPC + tar + gzip room on slow boxes without running forever.
ExecStart=/opt/blakestream-mpos/bin/backup.sh /var/backups/blakestream-mpos
StandardOutput=append:/var/log/blakestream-mpos/backup.log
StandardError=append:/var/log/blakestream-mpos/backup.log
TimeoutStartSec=15min
