# deploy-bundle/systemd/blakestream-mpos-backup.timer
[Unit]
Description=Schedule for Blakestream-MPOS backup snapshot
Requires=blakestream-mpos-backup.service

[Timer]
# Fire every 30 minutes. The actual decision to run lives in
# bin/backup.sh, which reads the operator-configured time from the
# settings table and bails cheaply if we're outside the window or a
# fresh backup already exists. This way the schedule is admin-editable
# without giving www-data write access to /etc/systemd/system/.
OnCalendar=*:0,30
RandomizedDelaySec=2m
Persistent=true
Unit=blakestream-mpos-backup.service

[Install]
WantedBy=timers.target
