# ops/systemd/mergeminer.service
# /etc/systemd/system/mergeminer.service
# Runs merged-mine-proxy-go against the parent + aux chain RPCs.
# sudo systemctl daemon-reload && sudo systemctl enable --now mergeminer

[Unit]
Description=Blakestream merge-mine proxy (auxpow relay for eloipool)
After=network.target eloipool.service
PartOf=eloipool.service

[Service]
Type=simple
User=eloipool
Group=eloipool
# Keep RPC credentials in the JSON config, not in process arguments. See
# ops/systemd/mergeminer.env.example.
Environment=MMP_WORKER_PORT=9992
EnvironmentFile=/etc/mpos/mergeminer.env
ExecStartPre=/usr/bin/test -s ${MMP_CONFIG}
ExecStart=/home/eloipool/eloipool-auxpow/bin/merged-mine-proxy-go --config ${MMP_CONFIG}
ExecStartPost=/home/eloipool/eloipool-auxpow/bin/check-mmp-ready.py ${MMP_WORKER_PORT} 5
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/blakecoin-pool/mmproxy.out
StandardError=append:/var/log/blakecoin-pool/mmproxy.err
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/var/log/blakecoin-pool
PrivateDevices=true
LockPersonality=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
